# scripts/ai.openclaw.tailscale-funnel.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
  <dict>
    <key>Label</key>
    <string>ai.openclaw.tailscale-funnel</string>

    <key>Comment</key>
    <string>Re-enable Tailscale Funnel after OpenClaw gateway restart (zaap-ho4)</string>

    <key>ProgramArguments</key>
    <array>
      <string>/bin/bash</string>
      <string>/Users/zane/Projects/zaap/scripts/ensure-funnel.sh</string>
    </array>

    <!-- Run at login/load and also watch the gateway log for restarts -->
    <key>RunAtLoad</key>
    <true/>

    <!-- Re-run whenever the gateway log is modified (gateway restarted) -->
    <key>WatchPaths</key>
    <array>
      <string>/Users/zane/.openclaw/logs/gateway.log</string>
    </array>

    <!-- Throttle: don't run more than once every 30 seconds -->
    <key>ThrottleInterval</key>
    <integer>30</integer>

    <key>StandardOutPath</key>
    <string>/Users/zane/.openclaw/logs/ensure-funnel.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/zane/.openclaw/logs/ensure-funnel.log</string>

    <key>EnvironmentVariables</key>
    <dict>
      <key>PATH</key>
      <string>/usr/local/bin:/usr/bin:/bin:/Applications/Tailscale.app/Contents/MacOS</string>
    </dict>
  </dict>
</plist>
